# deploy/systemd/bodydashboard-api.service
[Unit]
Description=Bodydashboard API
After=network.target

[Service]
Type=simple
User=bodydashboard
Group=bodydashboard
WorkingDirectory=/opt/bodydashboard/current/apps/api
EnvironmentFile=/etc/bodydashboard/bodydashboard.env
ExecStart=/usr/bin/node dist/src/server.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/bodydashboard /opt/bodydashboard

[Install]
WantedBy=multi-user.target
